# src/deb/init.d/oss-text-extractor
# oss-text-extractor start
#

### BEGIN INIT INFO
# Provides:          oss-text-extractor
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start/stop oss-text-extractor at boot time
# Description:       Enable service provided by oss-text-extractor.
### END INIT INFO

OSS_NAME=oss-text-extractor

CONFIGFILE=/etc/opensearchserver/$OSS_NAME
if [ -f "$CONFIGFILE" ]; then
    . "$CONFIGFILE"
fi

OSS_SHARE=/usr/share/opensearchserver
OSS_JAR=$OSS_SHARE/$OSS_NAME.jar
OSS_OPTS="$JAVA_OPTS -jar $OSS_JAR -p $SERVER_PORT"
OSS_LOG=/var/log/opensearchserver/$OSS_NAME/server.out
OSS_PID="/var/run/$OSS_NAME.pid"
OSS_USER=opensearchserver

usage ()
{
  echo "Usage : $0 {start|stop|restart}"
  exit
}

if [ "$#" -ne 1 ]
then
  usage
fi

case $1 in
start)
        if [ -f $OSS_PID ]; then
            echo "$OSS_NAME is already running ..."
        else
            echo "Starting $OSS_NAME ..."
            PID=`su -c "/usr/bin/java $OSS_OPTS" $OSS_USER >> $OSS_LOG 2>&1 & echo $!`
            if [ -z $PID ]; then
               echo "$OSS_NAME starting failed"
            else
               echo $PID > $SERVER_PID
               echo "$OSS_NAME started ..."
            fi
        fi
        ;;
stop)
        if [ -f $OSS_PID ]; then
            echo "Stopping $OSS_NAME ..."
            PID=`cat $OSS_PID`
            kill $PID
            echo "$OSS_NAME stopped ..."
            rm $OSS_PID
        else
            echo "$OSS_NAME is not running ..."
        fi
        ;;
status)
        if [ -f $OSS_PID ]; then
            PID=`cat $OSS_PID`
            echo "$OSS_NAME pid: $PID"
        else
            echo "$OSS_NAME is not running ..."
        fi
        ;;
restart|force-reload)
        $0 stop
        $0 start
		;;
esac
exit 0
