# src/deb/control/postinst
#!/bin/sh

if which invoke-rc.d >/dev/null 2>&1; then
	invoke-rc.d oss-text-extractor start
else
	/etc/init.d/oss-text-extractor start
fi
